# atxraspi-shutdown.service
# ATX-Raspi / systemd unit for Debian Bookworm (no rc.local)
# Runs the interrupt-driven shutdown/power watcher (shutdownirq.py)

[Unit]
Description=ATX-Raspi shutdown/power watcher (interrupt-driven)
After=multi-user.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /usr/local/sbin/atxraspi-shutdownirq.py
Restart=always
RestartSec=1
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
